# univision2022backend.service
[Service]
ExecStart=/home/www/univision2022backend/start.sh
Restart=always

[Install]
WantedBy=multi-user.target
